# com.aupat.worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aupat.worker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/bryant/Documents/tools/aupat/scripts/archive_worker.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/bryant/Documents/tools/aupat</string>

    <key>StandardOutPath</key>
    <string>/Users/bryant/Documents/tools/aupat/logs/worker.out</string>

    <key>StandardErrorPath</key>
    <string>/Users/bryant/Documents/tools/aupat/logs/worker.err</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
